# crontabs/hexagon_crontab
# DO NOT EDIT THIS FILE - edit the master and reinstall.
# (/tmp/crontab.XXXXO0wNqG installed on Fri Jul 24 14:57:36 2015)
# (Cron version V5.0 -- $Id: crontab.c,v 1.12 2004/01/23 18:56:42 vixie Exp $)
# Crontabs of Tim on Hexagon run on node1
# To set up cron job on ONE node, a) > ssh login1, > b) crontab hexagon_crontab, c) > exit 
#

# MAILTO="user@domain.com"
# MAILTO="timothy.williams@nersc.no"

# set variables
# * NB1 : these are "global variables" - all scripts launched from here will be able to access them
# * NB2 : scripts launched from crontab don't have access to .bash_profile/rc variables
SWARP_ROUTINES="/home/nersc/timill/GITHUB-REPOSITORIES/SWARP-routines"
FORECAST="/home/nersc/timill/GITHUB-REPOSITORIES/SWARP-routines/forecast_scripts"
VALIDATION="/home/nersc/timill/GITHUB-REPOSITORIES/SWARP-routines/validation"
qsub="/opt/torque/5.1.1-1_fba25d92/bin/qsub"
python="/work/apps/python/2.7.9-dso-gnu/bin/python"
qstat="/opt/torque/5.1.1-1_fba25d92/bin/qstat"
cdo="/work/apps/cdo/1.5.4-cray/bin/cdo"

# get wamnsea waves as early as possible
# NOTE the wamnsea products are usually uploaded just after 4 am
00 * * * * $FORECAST/wamnsea_update.sh

# TP4 restarts less urgent
00 18 * * * $FORECAST/topaz_archive_restart.sh

# run ice_only forecast as early as possible
00 02 * * * $FORECAST/ice_only/run_forecast.sh

# run wavesice forecast IF it didn't run after ice_only
# NOTE it runs after ice_only only if it finds the waves ice hence never because waves products are uploaded at 4 am
10 * * * * $FORECAST/wavesice/run_forecast_wav.sh

# get ice conc for validation and wave warning
00 12 * * * $VALIDATION/get_osisaf_conc.sh

# run check_wamsea.py (alert if big waves near ice)
10 ** * * * $FORECAST/alert_waves/check_wamnsea.sh
