# crontabs/johansen_crontab
# Crontab for SWARP project on Johansen
# 1) Getting the final product from Hexagon - try every hour

00 * * * * /Data/sim/tim/Projects/SWARP/SWARP-routines/forecast_scripts/copy2johansen.sh
