# crontabs/timcrontab
# DO NOT EDIT THIS FILE - edit the master and reinstall.
# (/tmp/crontab.XXXXiEFJRS installed on Mon Apr 20 10:05:52 2015)
# (Cron version V5.0 -- $Id: crontab.c,v 1.12 2004/01/23 18:56:42 vixie Exp $)
# DO NOT EDIT THIS FILE - edit the master and reinstall.
# (/tmp/crontab.XXXX3F6iX1 installed on Fri Apr 10 11:18:27 2015)
# (Cron version V5.0 -- $Id: crontab.c,v 1.12 2004/01/23 18:56:42 vixie Exp $)
# Crontabs of Tim on Hexagon run on node1
# To set up cron job on ONE node, a) > ssh login1, > b) crontab timcrontab, c) > exit 
# 1) Check restarts from topaz and wamnsea               - Every day @18.00
# 2) Forecast run                                        - Every day @02.00 
# 3) Download validation data from OSI_SAF               - Every day @12.00 
#
#MAILTO="gcmdnt90@gmail.com"
FORECAST="/home/nersc/timill/GITHUB-REPOSITORIES/SWARP-routines/forecast_scripts"
VALIDATION="/home/nersc/timill/GITHUB-REPOSITORIES/SWARP-routines/validation"

# get wamnsea waves as early as possible
00 * * * * $FORECAST/wamnsea_update.sh

# TP4 restarts less urgent
00 18 * * * $FORECAST/topaz_archive_restart.sh

# run forecast as early as possible
00 02 * * * $FORECAST/run_forecast.sh

# get ice conc for validation and wave warning
00 12 * * * $VALIDATION/get_osisaf_conc.sh
