# crontabs/timcrontab
# Crontabs of Tim on Hexagon run on node1
# To set up cron job on ONE node, a) > ssh login1, > b) crontab timcrontab, c) > exit 
# 1) Check restarts from topaz                           - Every day @1800
# 2) Check wamsea                                        - Every day @1800
# 2) Forecast run                                        - Every day @1900 
# 3) Archiviation of model's products                    - Every day @2100
# 4) 
#
MAILTO="timothy.williams@nersc.no,gcmdnt90@gmail.com"
SWARP_ROUTINES="/home/nersc/timill/GITHUB-REPOSITORIES/SWARP-routines"
TP4_REALTIME="/work/timill/RealTime_Models/TP4a0.12"

00 18 * * * $SWARP_ROUTINES/forecast_scripts/topaz_archive_restart.sh
00 18 * * * $SWARP_ROUTINES/forecast_scripts/wamnsea_update.sh
00 19 * * * $SWARP_ROUTINES/forecast_scripts/run_forecast.sh
00 21 * * * $SWARP_ROUTINES/forecast_scripts/process_FCresults.sh

