# crontabs/timcrontab
# DO NOT EDIT THIS FILE - edit the master and reinstall.
# (/tmp/crontab.XXXXyFINWz installed on Tue Jul 14 18:26:46 2015)
# (Cron version V5.0 -- $Id: crontab.c,v 1.12 2004/01/23 18:56:42 vixie Exp $)
# DO NOT EDIT THIS FILE - edit the master and reinstall.
# (/tmp/crontab.XXXXG5TyDu installed on Sun Jun 28 16:58:34 2015)
# (Cron version V5.0 -- $Id: crontab.c,v 1.12 2004/01/23 18:56:42 vixie Exp $)
# DO NOT EDIT THIS FILE - edit the master and reinstall.
# (/tmp/crontab.XXXXd4LMz2 installed on Wed Jun 24 15:08:02 2015)
# (Cron version V5.0 -- $Id: crontab.c,v 1.12 2004/01/23 18:56:42 vixie Exp $)
# DO NOT EDIT THIS FILE - edit the master and reinstall.
# (/tmp/crontab.XXXXiEFJRS installed on Mon Apr 20 10:05:52 2015)
# (Cron version V5.0 -- $Id: crontab.c,v 1.12 2004/01/23 18:56:42 vixie Exp $)
# DO NOT EDIT THIS FILE - edit the master and reinstall.
# (/tmp/crontab.XXXX3F6iX1 installed on Fri Apr 10 11:18:27 2015)
# (Cron version V5.0 -- $Id: crontab.c,v 1.12 2004/01/23 18:56:42 vixie Exp $)
# Crontabs of Tim on Hexagon run on node1
# To set up cron job on ONE node, a) > ssh login1, > b) crontab timcrontab, c) > exit 
#

#MAILTO="user@domain.com"
#MAILTO="timothy.williams@nersc.no"

# set variables
# * NB1 : these are "global variables" - all scripts launched from here will be able to access them
# * NB2 : scripts launched from crontab don't have access to .bash_profile/rc variables
FORECAST="/home/nersc/timill/GITHUB-REPOSITORIES/SWARP-routines/forecast_scripts"
VALIDATION="/home/nersc/timill/GITHUB-REPOSITORIES/SWARP-routines/validation"
qsub="/opt/torque/5.1.1-1_fba25d92/bin/qsub"
python="/work/apps/python/2.7.9-dso-gnu/bin/python"

# get wamnsea waves as early as possible
# NOTE the wamnsea products are uploaded AFTER 4 am
# added another download at 5 am
00 * * * * $FORECAST/wamnsea_update.sh
05 * * * * $FORECAST/wamnsea_update.sh

# TP4 restarts less urgent
00 18 * * * $FORECAST/topaz_archive_restart.sh

# run ice_only forecast as early as possible
00 02 * * * $FORECAST/ice_only/run_forecast.sh

# run wavesice forecast IF it didn't run after ice_only
# NOTE it runs after ice_only only if it finds the waves ice hence never because waves products are uploaded at 4 am
06 * * * * $FORECAST/wavesice/run_forecast_wav.sh

# get ice conc for validation and wave warning
00 12 * * * $VALIDATION/get_osisaf_conc.sh

# run check_wamsea.py
10 12 * * * $python $FORECAST/alert_waves/check_wamnsea.py

