# crontabs/timcrontab
# Crontabs of Tim on Hexagon run on node1
# To set up cron job on ONE node, a) > ssh login1, > b) crontab timcrontab, c) > exit 
# 1) Check restarts from topaz and wamnsea               - Every day @18.00
# 2) Forecast run                                        - Every day @19.00 
#
MAILTO="timothy.williams@nersc.no,gcmdnt90@gmail.com"
FORECAST="/home/nersc/timill/GITHUB-REPOSITORIES/SWARP-routines/forecast_scripts"

00 18 * * * $FORECAST/archive_inputs.sh
00 19 * * * $FORECAST/run_forecast.sh

